# btrfs-usb-backup/systemd/usb-backup-reminder.timer
[Unit]
Description=usb-backup-reminder.timer

[Timer]
OnCalendar=Fri *-*-* 17:00:00
Persistent=true
Unit=usb-backup-reminder.service

[Install]
WantedBy=timers.target
